# devops/crontab
*/5 * * * * sleep $((RANDOM\%15)) && ionice -c idle nice -19 /bin/bash /home/atomicpi/dns_warmer/devops/dns_warmer.sh  > /dev/null 2>&1
