# devops/crontab
*/5 * * * * /bin/bash /home/atomicpi/dns_warmer/devops/dns_warmer.sh
